# linux-gnu/etc/systemd/system/zram_swap.service.d/override.conf
[Service]
ExecStart=/bin/sh -c "exec /sbin/zram-init -azstd -Lzram_swap `LC_ALL=C free -m | awk '/^Mem:/{print int($2/4)}'`"
